# install/notes-db-backup.service
[Service]
Type=oneshot
ExecStart=/home/ubuntu/notes-service-install/backup-db.sh
